# templates/default/zone_profile.xml.erb
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="profile" name="system configuration">
  <service name="system/environment" version="1" type="service">
    <instance name="init" enabled="true">
      <property_group name="environment" type="application">
        <propval name="LANG" type="astring" value="en_US.UTF-8"/>
      </property_group>
    </instance>
  </service>
  <service name="system/console-login" version="1" type="service">
    <instance name="default" enabled="true">
      <property_group name="ttymon" type="application">
        <propval name="terminal_type" type="astring" value="xterm"/>
      </property_group>
    </instance>
  </service>
  <service name="system/timezone" version="1" type="service">
    <instance name="default" enabled="true">
      <property_group name="timezone" type="application">
        <propval name="localtime" type="astring" value="US/Pacific"/>
      </property_group>
    </instance>
  </service>
  <service name="system/config-user" version="1" type="service">
    <instance name="default" enabled="true">
      <property_group name="root_account" type="application">
        <propval name="login" type="astring" value="root"/>
        <propval name="password" type="astring" value="<%= @password_hash %>"/>
        <propval name="type" type="astring" value="role"/>
      </property_group>
    </instance>
  </service>
  <service name="system/identity" version="1" type="service">
    <instance name="node" enabled="true">
      <property_group name="config" type="application">
        <propval name="nodename" type="astring" value="<%= @hostname %>"/>
      </property_group>
    </instance>
  </service>
  <service name="network/physical" version="1" type="service">
    <instance name="default" enabled="true">
      <property_group name="netcfg" type="application">
        <propval name="active_ncp" type="astring" value="DefaultFixed"/>
      </property_group>
    </instance>
  </service>
  <service name="network/ssh" version="1" type="service">
    <instance enabled="true" name="default"/>
  </service>
  <service version="1" type="service" name="system/name-service/cache">
    <instance enabled="false" name="default"/>
  </service>
  <service version="1" type="service" name="network/smtp">
    <instance enabled="false" name="default"/>
  </service>
  <service version="1" type="service" name="network/dns/client">
    <instance enabled="true" name="default"/>
  </service>
</service_bundle>
